# supervisor.conf
[supervisord]
nodaemon=true
user=looker

[include]
files = /etc/supervisor/conf.d/*.conf

[program:looker]
command=/home/looker/looker/looker.sh start
user=looker